# node-js/dockerfile/start-service.ini
[supervisord]
nodaemon=true

[program:start_node_service]
command = /usr/local/bin/npm start
directory = /opt/gitalb_jenkins_lab
user = root
stopsignal = INT  
autostart = true
startsecs = 1
